# deploy/oracle/research-gap-bot.service
[Unit]
Description=Research Gap Analyzer Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/research-gap-analyzer
EnvironmentFile=/opt/research-gap-analyzer/.env
ExecStart=/opt/research-gap-analyzer/.venv/bin/python /opt/research-gap-analyzer/telegram_bot.py
Restart=always
RestartSec=10
KillSignal=SIGINT
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
